# beast.service.example
[Unit]
Description=Beast server
After=network.target

[Service]
User=vagrant
ExecStart=/home/vagrant/go/bin/beast run -v

[Install]
WantedBy=multi-user.target
